# supervisord.conf
[supervisord]
nodaemon=true

[program:sickbeard]
command=python /Sick-Beard/SickBeard.py
